# academic_radar/com.zhaomuchuan.academic-radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zhaomuchuan.academic-radar</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/zhaomuchuan/Documents/ai使用研究/academic_radar/scripts/run_daily.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <key>LC_ALL</key>
    <string>en_US.UTF-8</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/zhaomuchuan/Documents/ai使用研究/academic_radar/output/logs/launchd_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/zhaomuchuan/Documents/ai使用研究/academic_radar/output/logs/launchd_stderr.log</string>
</dict>
</plist>
